# rsync_backup.timer
[Unit]
Description=Weekly rsync backup

[Timer]
OnCalendar=Sun *-*-* 23:50

[Install]
WantedBy=timers.target
